# scripts/systemd/qwen36-multi.service
# Deprecated compatibility unit for the old qwen36-multi profile.
#
# Upstream llama.cpp now supports MTP with --parallel >1, and qwen36.service
# defaults to PARALLEL=4. This unit remains only so old installs keep working
# while they migrate.

[Unit]
Description=Deprecated Qwen3.6-27B compatibility shim (use qwen36.service with PARALLEL)
After=network-online.target
Conflicts=qwen36.service

[Service]
Type=simple
EnvironmentFile=%h/.config/qwen36-mtp/env
Environment=CUDA_VISIBLE_DEVICES=0
Environment=GGML_CUDA_ENABLE_UNIFIED_MEMORY=0
ExecStart=%h/Dev/qwen36-mtp/scripts/run-multi.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
